# new_android/app/config/port-forward.conf
[program:port-forward]
command=/app/scripts/port_forward.sh
autorestart=false
priority=1