# new_android/app/config/wmctrl.conf
[program:wmctrl]
command=wmctrl -c 'Running in Nested'
autorestart=false