# custom_android/app/config/adb-utils.conf
[program:adb-utils]
command=/app/scripts/utils.sh
autorestart=false
priority=5