# custom_android/app/config/auto-recording.conf
[program:auto-recording]
command=./app/scripts/record.sh auto_record
autorestart=false
priority=4