# custom_android/app/config/docker-appium.conf
[program:docker-appium]
command=/app/scripts/appium.sh
autorestart=true
priority=4